# deploy/roles/common/templates/supervisor.conf
[program:{{ app_name }}]
command=/usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/{{ app_name }}.ini
directory={{ app_path }}
stopsignal=QUIT
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/{{ app_name }}.log
environment=ATOLL_CONF={{ config_path }}
